# deploy/muer-assistant.service
[Unit]
Description=Muer Voice Assistant
Documentation=https://github.com/johnny4091/muer
After=sound.target pulseaudio.service pipewire.service
Wants=sound.target

[Service]
Type=simple
# Replace 'ubuntu' with your actual username
User=%i
Group=audio
WorkingDirectory=/opt/muer-assistant
EnvironmentFile=/opt/muer-assistant/.env
ExecStart=/opt/muer-assistant/venv/bin/python -m assistant.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Allow microphone access
SupplementaryGroups=audio pulse pulse-access

[Install]
WantedBy=multi-user.target
